# vermeer_inventory/vermeer_inventory.timer
[Unit]
Description=Generate and email inventory report to Vermeer - daily

[Timer]
OnCalendar=*-*-* 06:00:00

[Install]
WantedBy=timers.target
