# vermeer_inventory/vermeer_inventory.service
[Unit]
Description=Generate and email inventory report to Vermeer

[Service]
Type=oneshot
ExecStart=/path/to/vermeer_inventory.rb
